# wdbx-toolbox.service
[Unit]
Description=WDBX Radio Toolbox
Documentation=https://github.com/wdbx/toolbox
After=network-online.target
Wants=network-online.target
# Start after NAS automount if configured — comment out if not using automount
# After=mnt-wdbx\x2dshare.mount

[Service]
Type=simple
User=wdbx
WorkingDirectory=/home/wdbx/wdbx-toolbox

# Adjust python path if using a venv: /home/wdbx/wdbx-toolbox/venv/bin/python3
ExecStart=/usr/bin/python3 /home/wdbx/wdbx-toolbox/run.py

# Restart on crash, but not on clean exit (e.g. intentional shutdown)
Restart=on-failure
RestartSec=10
StartLimitIntervalSec=120
StartLimitBurst=5

# Logging — stdout/stderr go to journald; view with: journalctl -u wdbx-toolbox -f
StandardOutput=journal
StandardError=journal
SyslogIdentifier=wdbx-toolbox

# Environment
Environment=PYTHONUNBUFFERED=1
# Remove this line on the production box — reload is for local dev only
# Environment=WDBX_DEV=1

[Install]
WantedBy=multi-user.target
